# loterias-frontend-dev.service
[Unit]
Description=Loterias Analyzer Frontend (Development - Hot Reload)
After=network.target

[Service]
Type=simple
User=vunter
WorkingDirectory=/home/vunter/projects/loterias/frontend
Environment=PATH=/home/vunter/.local/share/fnm/node-versions/v24.13.0/installation/bin:/usr/bin:/bin
Environment=NODE_ENV=development
Environment=BACKEND_URL=http://127.0.0.1:8080
ExecStart=/home/vunter/.local/share/fnm/node-versions/v24.13.0/installation/bin/node /home/vunter/.local/share/fnm/node-versions/v24.13.0/installation/bin/npm run dev
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
