# loterias-backend-dev.service
[Unit]
Description=Loterias Analyzer Backend (Development - Hot Reload)
After=network.target

[Service]
Type=simple
User=vunter
WorkingDirectory=/home/vunter/projects/loterias/backend
EnvironmentFile=/home/vunter/projects/loterias/loterias-env.conf
ExecStart=/bin/bash -c 'source /home/vunter/projects/loterias/backend/vault-env.sh && exec /home/vunter/projects/loterias/backend/mvnw spring-boot:run -Dspring-boot.run.jvmArguments="-Dspring.devtools.restart.enabled=true"'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=SPRING_PROFILES_ACTIVE=dev,datadog
TimeoutStartSec=120
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM
SuccessExitStatus=143
MemoryMax=1G
CPUQuota=200%
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
