# loterias-backend.service
[Unit]
Description=Loterias Analyzer Backend (Spring Boot)
After=network.target

[Service]
Type=simple
User=loterias
Group=loterias
WorkingDirectory=/home/vunter/projects/loterias/backend
EnvironmentFile=/home/vunter/projects/loterias/loterias-env.conf
ExecStart=/bin/bash -c 'source /home/vunter/projects/loterias/backend/vault-env.sh && exec /usr/bin/java -javaagent:/home/vunter/projects/loterias/backend/lib/dd-java-agent.jar -Ddd.service=loterias -Ddd.env=prod -Ddd.version=0.0.1 -Ddd.logs.injection=true -Ddd.trace.sample.rate=0.5 -Ddd.profiling.enabled=true -Ddd.appsec.enabled=false -jar target/loterias-analyzer-0.0.1-SNAPSHOT.jar --server.address=127.0.0.1 --server.port=8081'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=SPRING_PROFILES_ACTIVE=prod,datadog
Environment=MANAGEMENT_DATADOG_METRICS_EXPORT_ENABLED=true
TimeoutStartSec=120
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM
SuccessExitStatus=143
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelTunables=true
ReadWritePaths=/home/vunter/projects/loterias/backend/target /var/log/loterias
MemoryMax=1G
CPUQuota=200%
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
