# loterias-frontend.service
[Unit]
Description=Loterias Analyzer Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=loterias
Group=loterias
WorkingDirectory=/home/vunter/projects/loterias/frontend
EnvironmentFile=/home/vunter/projects/loterias/loterias-env.conf
ExecStart=/home/vunter/.local/share/fnm/node-versions/v24.13.0/installation/bin/node /home/vunter/projects/loterias/frontend/.next/standalone/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=BACKEND_URL=http://127.0.0.1:8081
Environment=LOG_LEVEL=info
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelTunables=true
ReadWritePaths=/home/vunter/projects/loterias/frontend/.next
MemoryMax=512M
CPUQuota=100%

[Install]
WantedBy=multi-user.target
